# system-config/etc/systemd/system/zfs-scrub@.timer
[Unit]
Description=Scrub ZFS pool twice a month on %i

[Timer]
OnCalendar=*-*-1,15 2:00:00
Persistent=true

[Install]
WantedBy=timers.target

