# etc/systemd/system/home-root.mount
[Unit]

[Mount]
What=/root
Where=/home/root
Type=none
Options=bind

[Install]
WantedBy=multi-user.target

